# misc/infra/etc/systemd/system/link-friend.service
[Service]
Type=simple
ExecStart=/usr/local/bin/link-friend.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
