# misc/infra/var/lib/materialist/.config/systemd/user/spl-api.socket
[Socket]
ListenStream=127.0.0.1:28999

[Install]
WantedBy=default.target
