# misc/infra/var/lib/materialist/.config/systemd/user/steamcmd-image-pull.timer
[Timer]
OnCalendar=daily
RandomizedDelaySec=2hr
Persistent=true

[Install]
WantedBy=timers.target
