# misc/infra/var/lib/materialist/.config/systemd/user/steamcmd-image-pull.service
[Unit]
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/podman pull docker.io/steamcmd/steamcmd:alpine
